# rootfs/etc/services.d/hacs/run
#!/usr/bin/with-contenv bashio

hacs_components=/hacs/custom_components

cd /config

# 后台等待 HA 创建 configuration.yaml，然后追加 http 配置
(
  while [ ! -f /config/configuration.yaml ]; do
    sleep 5
  done
  # 额外等待，确保 HA 完成文件写入
  sleep 10
  if ! grep -q "^http:" /config/configuration.yaml; then
    bashio::log.info "Adding http configuration to configuration.yaml"
    cat >> /config/configuration.yaml <<EOF

http:
  use_x_forwarded_for: true
  trusted_proxies:
    - ::1
    - 127.0.0.1
    - 172.17.0.0/16
    - 172.18.0.0/16
    - 10.0.0.0/8
    - 192.168.0.0/16
EOF
  fi
) &

bashio::log.info "always_upgrade=$ALWAYS_UPGRADE"
if [ "$ALWAYS_UPGRADE" = "true" ]; then
    /install.sh
    ALWAYS_INSTALL=false
fi

[ ! -d custom_components ] && cp -rf $hacs_components ./ && bashio::log.info "Copied custom_components to config"

bashio::log.info "always_install=$ALWAYS_INSTALL"
if [ "$ALWAYS_INSTALL" = "true" ]; then
    cd custom_components
    [ ! -d hacs ] && cp -rf $hacs_components/hacs ./ && bashio::log.info "Copied hacs"
    [ ! -d xiaomi_miot ] && cp -rf $hacs_components/xiaomi_miot ./ && bashio::log.info "Copied xiaomi_miot"
    [ ! -d xiaomi_home ] && cp -rf $hacs_components/xiaomi_home ./ && bashio::log.info "Copied xiaomi_home"
    [ ! -d xiaomi_gateway3 ] && cp -rf $hacs_components/xiaomi_gateway3 ./ && bashio::log.info "Copied xiaomi_gateway3"
    [ ! -d sonoff ] && cp -rf $hacs_components/sonoff ./ && bashio::log.info "Copied sonoff"
    [ ! -d tianqi ] && cp -rf $hacs_components/tianqi ./ && bashio::log.info "Copied tianqi"
    [ ! -d edge_tts ] && cp -rf $hacs_components/edge_tts ./ && bashio::log.info "Copied edge_tts"
    [ ! -d ha_file_explorer ] && cp -rf $hacs_components/ha_file_explorer ./ && bashio::log.info "Copied ha_file_explorer"
    [ ! -d haier ] && cp -rf $hacs_components/haier ./ && bashio::log.info "Copied haier"
    [ ! -d midea_ac_lan ] && cp -rf $hacs_components/midea_ac_lan ./ && bashio::log.info "Copied midea_ac_lan"
    [ ! -d midea_ac ] && cp -rf $hacs_components/midea_ac ./ && bashio::log.info "Copied midea_ac"
    [ ! -d terncy ] && cp -rf $hacs_components/terncy ./ && bashio::log.info "Copied terncy"
    [ ! -d ai_conversation ] && cp -rf $hacs_components/ai_conversation ./ && bashio::log.info "Copied ai_conversation"
    [ ! -d zhixiangquanwu ] && cp -rf $hacs_components/zhixiangquanwu ./ && bashio::log.info "Copied zhixiangquanwu"
    cd ..
fi

sleep infinity